# ujobs/source/supervisord.conf
[program:agentinfo]
command={{ PYTHON }} {{ PROJECT_DIR }}/manage.py agentinfo