# monitor_dashboard/monitor_dashboard/supervisord.conf
[program:plugin_settings_cache]
command={{ PYTHON }} {{ PROJECT_DIR }}/manage.py update_cache